# config/php/supervisord.conf
[supervisord]
nodaemon=true

[program:laravel]
command=php artisan serve mylaravel --host 0.0.0.0

[program:thinkphp]
command=php artisan serve mythinkphp --host 0.0.0.0